# deployment/supervisor/aucr.conf
[program:aucr]
command=/home/ubuntu/aucr/venv/bin/gunicorn -b localhost:8000 -w 4 aucr:app
directory=/home/ubuntu/aucr
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
