# deployment/supervisor/aucr-tasks.conf
[program:aucr-tasks]
command=/home/ubuntu/aucr/venv/bin/rq worker aucr-tasks
numprocs=1
directory=/home/ubuntu/aucr
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
